# scripts/com.w3geekery.tt-dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.w3geekery.tt-web-servers</string>
    <key>Program</key>
    <string>/Users/cstacer/Projects/w3geekery/tt/scripts/tt-web-servers</string>
    <key>ProgramArguments</key>
    <array>
        <string>tt-web-servers</string>
        <string>start</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <!-- StartCalendarInterval removed 2026-04-16. The 6 AM trigger was
         firing during macOS DarkWake cycles, which then immediately killed
         the spawned processes on re-sleep. User-invoked start via
         `/tt hey` / `/tt good morning` is the explicit lifecycle instead.
         RunAtLoad covers the full-login case. -->
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>StandardOutPath</key>
    <string>/tmp/tt-dev.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tt-dev.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
